# security/hardening/build/debian-build/hardening-lsm-1.0.0/packaging/systemd/hardening-lsm.service
[Unit]
Description=Security Hardening LSM Management Service
After=multi-user.target
ConditionPathExists=/sys/kernel/security/hardening

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/hardening-ctl enable
ExecStop=/usr/bin/hardening-ctl disable
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target